# facebook-token-scheduler.service
[Unit]
Description=Facebook Token Refresh Scheduler
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/palimirror
ExecStart=/usr/bin/node facebook-token-scheduler.js
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=facebook-token-scheduler
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target 